# services/logrotate.timer
[Unit]
Description=Run logrotate periodically

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true

[Install]
WantedBy=timers.target
